# scripts/cron.d/htcondor-accounting
# Run the full daily pipeline once per day. The script computes yesterday in UTC by default.
15 01 * * * condor /path/to/htcondor-accounting/scripts/run_daily_pipeline.sh /srv/htcondor-accounting >> /var/log/htcondor-accounting-daily.log 2>&1
